# system_files/usr/lib/systemd/system/discord-ping.timer
[Unit]
Description=Run discord ping

[Timer]
OnBootSec=20s
RandomizedDelaySec=10s
OnUnitInactiveSec=2min

[Install]
WantedBy=timers.target
